# infra/bartleby.service
# Example systemd unit for running Bartleby as a container on bare metal.
# Adjust the image tag and paths, then:
#   sudo cp bartleby.service /etc/systemd/system/
#   sudo systemctl daemon-reload && sudo systemctl enable --now bartleby
#
# Expects an env file at /etc/bartleby/bartleby.env (see .env.example) and a
# vault directory at /var/lib/bartleby.

[Unit]
Description=Bartleby — self-hosted knowledge vault
After=network-online.target
Wants=network-online.target

[Service]
Restart=always
RestartSec=5
ExecStartPre=-/usr/bin/docker rm -f bartleby
ExecStart=/usr/bin/docker run --rm --name bartleby \
  --env-file /etc/bartleby/bartleby.env \
  -p 127.0.0.1:8080:8080 \
  -v /var/lib/bartleby:/data \
  ghcr.io/t11z/bartleby:latest
ExecStop=/usr/bin/docker stop bartleby

[Install]
WantedBy=multi-user.target
